# scripts/signal-noise.service
[Unit]
Description=signal-noise API
After=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/projects/signal-noise
ExecStart=/home/dev/projects/signal-noise/.venv/bin/python -m signal_noise serve --host 127.0.0.1 --port 8000 --no-scheduler
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
